# deployment/bulletin-server.service
[Unit]
Description=Bulletin Board Flask Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/bulletin
EnvironmentFile=/opt/bulletin/config.env
ExecStart=/opt/bulletin/venv/bin/python /opt/bulletin/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
