# deployment/kiosk-daemon.service
[Unit]
Description=Bulletin Board Kiosk Client
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=%I
Group=%I
WorkingDirectory=/opt/bulletin-client
EnvironmentFile=/opt/bulletin-client/config.env
ExecStart=/opt/bulletin-client/venv/bin/python /opt/bulletin-client/client_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
